# crontab
* * * * * /home/veronico/assign1/myspammer.sh

